# Library/LaunchDaemons/com.makkintosshu.tiggit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.makkintosshu.tiggit</string>
    <key>Program</key>
    <string>/usr/local/bin/tiggit</string>
    <key>ProgramArguments</key>
    <array>
      <string>tiggit</string>
      <string>--verbose</string>
      <string>--daemon</string>
    </array>
    <key>UserName</key>
    <string>_tiggit</string>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
  </dict>
</plist>
